# addons/host/etc/systemd/system/weston.service.d/lxc-net.conf
[Unit]
After=dbus.service multi-user.target lxc-net.service
