# addons/host/lib/systemd/system/lxc-watchdog-system.service
[Unit]
Description=AGL Bomb (IVI) 

[Service]
Type=simple
ExecStart=/usr/bin/watchdog2.sh
